# resource/supervisor/config/supervisor-service.conf
[program:cron]
command = /usr/sbin/cron -f -L 15
autostart = true
priority = 999