# docker/build/php7.3-nginx/config/php-fpm.conf
[program:php-fpm7.3]
command=/usr/sbin/php-fpm7.3 -F
autostart=true
autorestart=true
priority=999
username=www-data
autostart=true
autorestart=unexpected
redirect_stderr=true
exitcodes=0
