# ubuntu/build/nodejs8/config/nginx-supervisor.conf
[program:nginx]
command=/usr/sbin/nginx
priority=999
autostart=true
autorestart=unexpected
redirect_stderr=true
exitcodes=0
stdout_events_enabled=true
stderr_events_enabled=true