# ubuntu/build/nodejs8/config/php-fpm.conf
[program:php-fpm7.1]
command=/usr/sbin/php-fpm7.1 -F
autostart=true
autorestart=true
priority=999
username=www-data
autostart=true
autorestart=unexpected
redirect_stderr=true
exitcodes=0
